# dot_config/systemd/user/crush-signal-channel.service.tmpl
# Crush headless session (Signal channel arbiter) — tmux-managed user service.
# One persistent Crush process attached to the Signal MCP as a channel; it is
# the SOLE arbiter of messages that land on this machine's Signal account.
# systemd keeps it alive; tmux gives you a way in to watch/steer:
#
#   attach: crush-signal-channel attach     (or: tmux -L crush attach -t signal)
#   detach: ^b d                            (standard tmux keybind)
#   logs:   crush-signal-channel log        (or: journalctl --user -u crush-signal-channel -f)
#   enable: systemctl --user enable --now crush-signal-channel
#
# Requires the signal-cli daemon (JSON-RPC backend the signal MCP talks to) and
# vault-agent-rendered env (LLM provider keys — Crush needs at least one).
#
# A user service runs while you have a session. To keep it running when you're
# NOT logged in (unattended box, needs sudo ONCE):
#   sudo loginctl enable-linger {{ .chezmoi.username }}
# macOS counterpart: Library/LaunchAgents/rocks.stump.crush-signal-channel.plist.tmpl
[Unit]
Description=Crush headless (Signal channel arbiter, tmux-managed)
After=network-online.target vault-agent.service signal-daemon.service
Wants=network-online.target vault-agent.service
Requires=signal-daemon.service

[Service]
# tmux new-session -d double-forks, so let systemd track it as forking. There's
# no PID file — the ExecStop below tears the tmux session down cleanly, and the
# supervisor loop inside the pane keeps crush alive across TUI crashes.
Type=forking
RemainAfterExit=yes

# LLM keys (LITELLM_API_KEY / OPENAI_DIRECT_API_KEY / GEMINI_API_KEY / ZAI_API_TOKEN)
# reach Crush via vault-agent's rendered env. Without at least one, Crush's
# providers block is empty and it can't run.
EnvironmentFile=%h/.config/vault/secrets-static.env
Environment=HOME=%h
Environment=USER={{ .chezmoi.username }}
Environment=VAULT_ADDR={{ .vaultAddr }}

ExecStart=%h/.local/bin/crush-signal-channel.sh
ExecStop={{ lookPath "tmux" | default "/usr/bin/tmux" }} -L crush kill-session -t signal

Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
