# Library/LaunchAgents/rocks.stump.crush-signal-channel.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- macOS launchd counterpart of dot_config/systemd/user/crush-signal-channel.service.tmpl.
     Runs the tmux-wrapped headless Crush signal-channel arbiter as a per-user
     LaunchAgent so it comes up when you log in, stays up while the laptop is
     open, and pauses automatically while the machine sleeps.
       attach: crush-signal-channel attach     (tmux -L crush attach -t signal)
       logs:   crush-signal-channel log        (~/.local/share/crush-signal-channel/agent.log)
       load:   crush-signal-channel start      (launchctl load -w …)
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>rocks.stump.crush-signal-channel</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{ .chezmoi.homeDir }}/.local/bin/crush-signal-channel.sh</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <!-- launchd's PATH is minimal. Point it at brew (both prefixes) + system
             paths so the wrapper's lookPath resolves tmux / crush / uv / etc. -->
        <key>PATH</key>
        <string>{{ .chezmoi.homeDir }}/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>{{ .chezmoi.homeDir }}</string>
        <key>VAULT_ADDR</key>
        <string>{{ .vaultAddr }}</string>
    </dict>
    <!-- RunAtLoad brings it up on login; KeepAlive restarts if the wrapper (and
         therefore the tmux server) exits. Note: tmux new-session -d daemonizes,
         so the wrapper exits ~immediately with 0; KeepAlive re-runs it only if
         crush itself exits and the supervisor loop breaks out — the loop
         normally holds the pane forever. -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>StandardOutPath</key>
    <string>{{ .chezmoi.homeDir }}/.local/share/crush-signal-channel/agent.log</string>
    <key>StandardErrorPath</key>
    <string>{{ .chezmoi.homeDir }}/.local/share/crush-signal-channel/agent.log</string>
</dict>
</plist>
