# Library/LaunchAgents/rocks.stump.vault-agent-stale.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>rocks.stump.vault-agent-stale</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>{{ .chezmoi.homeDir }}/.config/vault/vault-agent-stale.sh</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>VAULT_ADDR</key>
        <string>{{ .vaultAddr }}</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>StandardOutPath</key>
    <string>{{ .chezmoi.homeDir }}/.config/vault/agent-stale.log</string>
    <key>StandardErrorPath</key>
    <string>{{ .chezmoi.homeDir }}/.config/vault/agent-stale.log</string>
</dict>
</plist>
