# Library/LaunchAgents/rocks.stump.vault-agent.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>rocks.stump.vault-agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{ (lookPath "vault") | default "/opt/homebrew/bin/vault" }}</string>
        <string>agent</string>
        <string>-config={{ .chezmoi.homeDir }}/.config/vault/agent.hcl</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>VAULT_ADDR</key>
        <string>{{ .vaultAddr }}</string>
        <!-- USER scopes the render to secret/users/$USER/*; HOME anchors
             ssh-keys.ctmpl's writeToFile targets at ~/.ssh. Both are read via env
             in the *.ctmpl files; launchd does not always populate them. -->
        <key>USER</key>
        <string>{{ .chezmoi.username }}</string>
        <key>HOME</key>
        <string>{{ .chezmoi.homeDir }}</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>{{ .chezmoi.homeDir }}/.config/vault/agent.log</string>
    <key>StandardErrorPath</key>
    <string>{{ .chezmoi.homeDir }}/.config/vault/agent.log</string>
</dict>
</plist>
